# bin/98-monitor-hotplug.rules
KERNEL=="card0", SUBSYSTEM=="drm", ACTION=="change", ENV{DISPLAY}=":0", ENV{XAUTHORITY}="/home/kareem/.Xauthority", RUN+="/bin/bash /home/kareem/bin/monitors"
